# supervisord/supervirsord.ini
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[program:mongod]
command=sh -c 'rm -rf /data/mongodb/mongod.lock && exec /usr/bin/mongod --port 27017 --dbpath /data/mongodb'
