# supervisord/supervirsord.conf
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[program:mongod]
command=/usr/bin/mongod --port 27017 --dbpath /data/mongodb
